# ubuntu1604/supervisord/conf.d/mysql.ini
[program:mysqld-chmod]
command=chown -R mysql:mysql /var/lib/mysql /var/run/mysqld

[program:mysqld]
command=/etc/init.d/mysql start
