# ubuntu1604/supervisord/conf.d/sshd.ini
[program:sshd]
priority = 5
command=/usr/sbin/sshd -D
autostart = true
autorestart = true
startsecs = 0
startretries = 0
redirect_stderr = true
stdout_logfile = /var/log/secure
stdout_events_enabled = true
